# systemd/dumpyarabot-bot.service
[Unit]
Description=Dumpyarabot Telegram bot
After=dumpyarabot-redis.service
Wants=dumpyarabot-redis.service
PartOf=dumpyarabot.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
User=jenkins
WorkingDirectory=/var/lib/jenkins/dumpbot
# pydantic-settings reads `.env` from the working directory, so no EnvironmentFile needed.
ExecStart=/var/lib/jenkins/dumpbot/.venv/bin/python -m dumpyarabot
Restart=on-failure
RestartSec=10
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
