# systemd/dumpyarabot-redis.service
[Unit]
Description=Redis for dumpyarabot (port 34790)
PartOf=dumpyarabot.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
User=jenkins
WorkingDirectory=/var/lib/jenkins/dumpbot
ExecStart=/usr/bin/redis-server /var/lib/jenkins/dumpbot/redis.conf --daemonize no
Restart=on-failure
RestartSec=5
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
