# systemd/dumpyarabot-worker@.service
[Unit]
Description=Dumpyarabot ARQ worker (instance %i)
After=dumpyarabot-redis.service
Wants=dumpyarabot-redis.service
PartOf=dumpyarabot.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
User=jenkins
# Worker shells out to uvx (~/.cargo/bin) and aria2c/7zz (~/.local/bin),
# neither of which is on the default system PATH.
Environment=PATH=/var/lib/jenkins/.cargo/bin:/var/lib/jenkins/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
WorkingDirectory=/var/lib/jenkins/dumpbot
ExecStart=/var/lib/jenkins/dumpbot/.venv/bin/python run_arq_worker.py worker_%i
Restart=on-failure
RestartSec=10
# Worker's job_timeout is 7200s (2h). Give it that plus a minute of grace
# so systemd lets an in-flight dump drain on SIGTERM instead of SIGKILLing
# mid-extraction.
TimeoutStopSec=7260
# If a child (aria2c, 7zz) is killed by the kernel OOM killer, don't
# mark the worker as failed — let it report the child error cleanly.
OOMPolicy=continue

[Install]
WantedBy=multi-user.target
