# systemd/dumpyarabot-worker.service
[Unit]
Description=Dumpyarabot ARQ worker
After=dumpyarabot-redis.service
Wants=dumpyarabot-redis.service
PartOf=dumpyarabot.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
WorkingDirectory=/var/lib/jenkins/dumpbot
# systemd user services start with a minimal PATH that excludes ~/.cargo/bin
# and ~/.local/bin. The worker shells out to `uvx`, `aria2c`, `7zz`, etc. —
# uvx in particular lives in ~/.cargo/bin/ for jenkins.
Environment=PATH=/var/lib/jenkins/.cargo/bin:/var/lib/jenkins/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/var/lib/jenkins/dumpbot/.venv/bin/python run_arq_worker.py
Restart=on-failure
RestartSec=10
KillSignal=SIGTERM
# Worker's job_timeout is 7200s (2h). Give it that plus a minute of grace
# so systemd lets an in-flight dump drain on SIGTERM instead of SIGKILLing
# mid-extraction.
TimeoutStopSec=7260

[Install]
WantedBy=default.target
